# systemd/netlink-dispatcher.service
[Service]
ExecStart=/usr/local/bin/netlink-dispatcher
